# rootfs/etc/service.d/whatwatt/run
#!/usr/bin/with-contenv bashio

# Pobranie konfiguracji
CONFIG_PATH="/data/options.json"
WHATWATT_IP=$(bashio::config 'whatwatt_ip')
MQTT_BROKER=$(bashio::config 'mqtt_broker')
MQTT_USERNAME=$(bashio::config 'mqtt_username')
MQTT_PASSWORD=$(bashio::config 'mqtt_password')
MQTT_PORT=$(bashio::config 'mqtt_port')
MQTT_TOPIC=$(bashio::config 'mqtt_topic')
REPORTING_INTERVAL=$(bashio::config 'reporting_interval')

# Uruchomienie skryptu konfiguracyjnego
exec python3 /usr/bin/whatwatt_integration.py \
    --whatwatt-ip "${WHATWATT_IP}" \
    --mqtt-broker "${MQTT_BROKER}" \
    --mqtt-username "${MQTT_USERNAME}" \
    --mqtt-password "${MQTT_PASSWORD}" \
    --mqtt-port "${MQTT_PORT}" \
    --mqtt-topic "${MQTT_TOPIC}" \
    --reporting-interval "${REPORTING_INTERVAL}" \
    --config "${CONFIG_PATH}"